# etc/profile.d/lang.sh
#!/bin/sh

source /etc/locale.conf
export LC_ALL=$LANG
